# .config/supervisord.ini
[program:horizon]
process_name=%(program_name)s
command=php /var/www/html/project/artisan horizon
autostart=true
autorestart=true
user=nginx
redirect_stderr=true
stdout_logfile=/var/log/supervisor/project-horizon.log
stopwaitsecs=3600